# chef/cookbooks/corosync/templates/default/corosync-shutdown-cleaner.service.erb
[Unit]
Description=Clean file blocking <%= @service_name %> start after fencing
Requires=<%= @service_name %>.service
Before=<%= @service_name %>.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/true
ExecStop=/usr/bin/rm -f <%= @block_corosync_file %>
ExecStopPost=/usr/bin/echo "corosync-shutdown-cleaner.service cleaning up <%= @block_corosync_file %> done"

[Install]
WantedBy=multi-user.target
