# chef/cookbooks/corosync/templates/default/corosync-shutdown-cleaner.service.erb
[Unit]
Description=Clean file blocking <%= @service_name %> start after fencing
Before=<%= @service_name %>.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=
ExecStop=rm -f <%= @block_corosync_file %>

[Install]
WantedBy=multi-user.target
